# deploy/templates/netcheck.timer
[Unit]
Description=Run netcheck regularly

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
